# rootfs/etc/cron.d/mopidy
* * * * * mediadepot s6-setuidgid mediadepot mopidy --config /srv/mopidy/config/mopidy.conf local scan
# An empty line is required at the end of this file for a valid cron file.
